# docker/crontab
# generated by ai

# Ensure correct PATH for npm and other binaries in cron's environment
# /usr/local/bin is where npm executables are often found for Node.js installed via official methods
# /usr/bin and /bin are for common Linux binaries
PATH=/usr/local/bin:/usr/bin:/bin

# Run the 'process-news' script at 9am and 9pm latvian time every day
# The timezone is set to Europe/Riga, which is the timezone for Latvia
0 9,21 * * * /usr/bin/env sh -c "cd /app && npm run process-news" > /dev/stdout 2> /dev/stderr

# An empty line is required at the end of a crontab file for POSIX compliance